# packaging/etc/supervisord.d/01-multiworker.conf
[program:multiworker]
command=bash -c "MultiWorker --action start"
priority=-1
stdout_logfile=/opt/snmpmon/multiworker-daemon.log
stdout_logfile_maxbytes=0
redirect_stderr=true
